# usr/lib/udev/rules.d/90-intel-rapl.rules
SUBSYSTEM=="powercap", KERNEL=="intel-rapl:*", ATTR{energy_uj}!="", RUN+="/bin/chmod 0444 /sys/class/powercap/%k/energy_uj"